# etoyataxi_debian/etoyataxi_api_stable/debian/ostdlabs-etoyataxi-api.cron.d
0 3 * * * etoyataxi /usr/share/ostdlabs/etoyataxi/etoyataxi-api/sql/control/backup.sh
0 1 */8 * * etoyataxi /usr/share/ostdlabs/etoyataxi/etoyataxi-api/sql/control/delete_old_backup.sh

